# ansible/roles/api/files/vAPI.conf
[program: vAPI]

command=/usr/local/bin/python /var/www/api/vAPI.py

directory=/var/www/api/

username=root

autostart=true

autorestart=true

redirect_stderr=true
